# ansible/roles/ilo/templates/fancontrol-ilo.timer.j2
# {{ ansible_managed }}

[Unit]
Description=iLO Fan Speed Control

[Timer]
OnBootSec=15sec
OnUnitActiveSec={{ ilo__run_interval }}

[Install]
WantedBy=timers.target
