# ansible/roles/ilo/templates/fancontrol-ilo.service.j2
# {{ ansible_managed }}

[Unit]
Description=iLO Fan Speed Control

[Service]
Type=oneshot
User=root
Group=root
ExecStart={{ ilo__fancontrol_executable }}
SyslogIdentifier=fancontrol-ilo
